# platform/linux/systemd/claude-metrics-ingest.service
[Unit]
Description=Rebuild Claude Code token/cost history (claude-metrics)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/claude-metrics ingest
